# etc/crontab
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=root
HOME=/

@daily root /usr/local/acme.sh/acme.sh --renew -d tjdc.co.uk && service nginx restart

